# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:prometheus]
command=/opt/prometheus/prometheus \
    --config.file=/etc/prometheus/prometheus.yml \
    --storage.tsdb.path=/prometheus \
    --storage.tsdb.retention.time=7d \
    --storage.tsdb.retention.size=256MB \
    --web.listen-address=:9090 \
    --web.enable-lifecycle
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=1

[program:grafana]
command=/opt/grafana/bin/grafana-server \
    --homepath=/opt/grafana \
    --config=/etc/grafana/grafana.ini \
    web
environment=GF_PATHS_DATA="/var/lib/grafana",GF_PATHS_LOGS="/var/log/grafana",GF_PATHS_PLUGINS="/var/lib/grafana/plugins",GF_PATHS_PROVISIONING="/etc/grafana/provisioning",GF_SERVER_HTTP_PORT="3000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=2
